# dockerfiles/ai-workspace/supervisord.conf
[supervisord]
# Run as the hermes user (set in the main config, but child procs inherit)
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp/supervisor-logs

# Avoid /var/run permission issues — use /tmp for socket
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ── CodeNomad Server ──────────────────────────────────────────────
# Serves the web UI for OpenCode sessions on port 9899
[program:codenomad]
command=codenomad --https=false --http=true --host 0.0.0.0 --workspace-root /opt/data --unrestricted-root --dangerously-skip-auth
directory=/opt/data
user=hermes
environment=HOME="/opt/data/home",NODE_ENV="production"
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Hermes Gateway ─────────────────────────────────────────────────
# Connects to Telegram (and other configured platforms) for agent access
[program:hermes-gateway]
command=/opt/hermes/.venv/bin/hermes gateway run
user=hermes
environment=HERMES_HOME="/opt/data",HOME="/opt/data/home",PATH="/opt/hermes/.venv/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
startretries=5
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
